# crontab
*/15 * * * * cd /app && ./docker-entrypoint.sh 1>/proc/1/fd/1 2>/proc/1/fd/2
